# etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=/sbin/agetty --skip-login --noissue --login-options "-f lerax" %I $TERM
